# extras/rtl_tcp.service
# /etc/systemd/system/rtl_tcp.service
[Unit]
Description=rtl_tcp SDR server
After=network.target

[Service]
ExecStart=/usr/local/bin/rtl_tcp -a 127.0.0.1 -p 1234
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
